# src/src/main/packages/deb/control-deb14/postinst
#!/bin/sh

invoke-rc.d rest-adapter-service start
