# src/src/main/packages/systemd/rest-adapter-service.service
[Unit]
Description=REST Adapter Service
After=network.target

[Service]
PermissionsStartOnly=true
User=rest-adapter-service
Group=rest-adapter-service
ExecStart=/usr/bin/rest-adapter-service
Restart=on-failure

[Install]
WantedBy=multi-user.target

