# src/src/main/packages/deb/control-deb16/postinst
#!/bin/sh

systemctl preset rest-adapter-service
systemctl start rest-adapter-service